# packages/lfrppy/lfrppy-0.0.20.tar.gz/lfrppy-0.0.20/lfrppy/resources/frpc/systemd/frpc@.service
[Unit]
Description=FRP Client Service (%i)
After=network.target
Wants=network.target

[Service]
WorkingDirectory=<dir>
Type=simple
User=<user>
Restart=on-failure
RestartSec=5s
ExecStart=<bin_dir>/frpc -c <dir>/%i.toml
ExecReload=<bin_dir>/frpc reload -c <dir>/%i.toml
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
